# contrib/certbot-dns-gandi-renew.timer
[Unit]
Description=Certbot weekly renewal

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
